# service/plex-sleep-guardian.service
[Unit]
Description=Plex: keep the machine from sleeping while streams are active
After=network-online.target
Wants=network-online.target

[Service]
Type=simple                # executa continuamente
ExecStart=/usr/local/bin/plex-sleep-guardian.sh
Restart=always

[Install]
WantedBy=multi-user.target
